# rest/skyght.service
[Unit]
Description=Skyght Business Flask Application
After=network.target postgresql.service
Wants=postgresql.service

[Service]
User=administrator
Group=administrator
WorkingDirectory=/opt/rest/skyght_business
Environment="PATH=/home/administrator/.local/bin:/usr/local/bin:/usr/bin"
EnvironmentFile=/opt/rest/skyght_business/.env
ExecStart=/home/administrator/.local/bin/gunicorn --workers 4 --bind 0.0.0.0:8080 --access-logfile - --error-logfile - "app:create_app()"
ExecReload=/bin/kill -s HUP $MAINPID
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
